# supervisord.conf
[supervisord]
nodaemon=true

[program:bind9]
command=/usr/sbin/named -f -u bind
